# files_rpi/etc/rc.local
#!/bin/sh -e

chmod +x /home/pi/.config/openbox/autostart /home/pi/tty2rpi-scripts/*
ln -sf /bin/bash /bin/sh
exit 0
